# crontabs/root.sample
55 21 * * * ./bin/disable-policy.sh projects/<Project ID>/alertPolicies/<Policy ID>
00 08 * * * ./bin/enable-pilicy.sh projects/<Project ID>/alertPolicies/<Policy ID>
